# iot-node/iot-node.service
[Unit]
Description=Node of HomeMonitor project

[Service]
ExecStart=/share/homemonitor/iot-node.sh

[Install]
WantedBy=multi-user.target


